# server/config/gunicorn.conf
description "Enklave Api "

start on runlevel [2345]
stop on runlevel [!2345]

respawn 
setuid ubuntu
setgid www-data

env PATH=/srv/enklave_app/env/bin
chdir /srv/enklave_app/enklave
exec gunicorn --workers 5 --bind unix:/srv/gunicorn/gunicorn.sock -m 007 enklave.wsgi:application >>/dev/null