# server/config/rqworker.conf
# Ubuntu upstart file at /etc/init/<service>.conf
# More info: http://upstart.ubuntu.com/cookbook/

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

setuid ubuntu
setgid ubuntu

exec /srv/hyde_api2_app/env/bin/python /srv/hyde_api2_app/hyde-backend/manage.py rqworker default